# release/scripts/rpm/websocketd.service
[Unit]
Description=websocketd daemon

[Service]
Type=simple
ExecStart=/usr/libexec/run-websocketd

[Install]
WantedBy=multi-user.target
